# autostart/.config/autostart/fix-hdmi-monitor.desktop
[Desktop Entry]
Type=Application
Name=Fix HDMI Monitor Resolution
Exec=/home/dblanken/.local/bin/enable-hdmi.sh
X-GNOME-Autostart-enabled=true
