# autostart/.config/autostart/vrr-enable.desktop
[Desktop Entry]
Type=Application
Name=Enable VRR/FreeSync
Exec=/home/dblanken/.local/bin/enable-vrr.sh
X-GNOME-Autostart-enabled=true
